# share/supervisor/sv.ini
[program:texas_web]
environment=PYTHON_EGG_CACHE=/tmp/.python-eggs/,MODE=CN
directory=/home/gengleiming/gengleiming
command=/usr/local/bin/gunicorn -c gunicorn_config.py gengleiming.wsgi:application
user=gengleiming
autorestart=true
redirect_stderr=true
