# share/supervisor/sv.ini
[program:gengleiming_web]
directory=/home/gengleiming/gengleiming
command=/usr/local/bin/gunicorn -c gunicorn_config.py gengleiming.wsgi:application
user=gengleiming
autorestart=true
redirect_stderr=true

[program:mkdocs]
directory=/home/gengleiming/mkcloud
command=mkdocs serve
user=gengleiming
autorestart=true
redirect_stderr=true
