# share/supervisor/sv.ini
[program:gengleiming_web]
directory=/data/release/gengleiming
command=gunicorn -c gunicorn_config.py gengleiming.wsgi:application
user=gengleiming
autorestart=true
redirect_stderr=true

[program:mkdocs]
directory=/data/release/mkcloud
command=mkdocs serve
user=gengleiming
autorestart=true
redirect_stderr=true
